# systemd/system/autoupload-rclone@.path
[Unit]
Description=Watch the %f directory for changes
After=tailscale-online@meerkat.target
BindsTo=tailscale-online@meerkat.target
StartLimitIntervalSec=0

[Path]
PathChanged=%f
MakeDirectory=true
TriggerLimitIntervalSec=0

[Install]
UpheldBy=tailscale-online@meerkat.target
WantedBy=multi-user.target
