# systemd/user/autoupload-rclone@.path
[Unit]
Description=Watch the %f directory for changes

[Path]
PathChanged=%f
MakeDirectory=true
TriggerLimitBurst=0
TriggerLimitIntervalSec=2min

[Install]
WantedBy=default.target
