# systemd/user/autoupload-immich.path
[Unit]
Description=Watch the directory for changes
After=tailscale-online@meerkat.target
BindsTo=tailscale-online@meerkat.target
StartLimitIntervalSec=0

[Path]
# todo Make this an instantiable unit when only versions five and newer of Podman are supported.
# PathChanged=%f
PathChanged=/home/jordan/Pictures
MakeDirectory=true
TriggerLimitIntervalSec=0

[Install]
WantedBy=default.target
